# projects/_template/app.container
# DEPRECATED — this stack has migrated from rootless Podman Quadlets to Docker Compose.
# Use compose.yaml in this directory instead.
#
# This file is retained for historical reference only. Podman Quadlet support
# may be revisited if NemoClaw adds a Podman driver — see docs/current/todos.md.
#
# ---- original template below (do not use) ------------------------------------
#
# Quadlet template (systemd --user).
# Copy this dir:   cp -r projects/_template projects/<name>
# Rename file to:  <name>.container   and replace every <name> below.
# Install:         ln -sf ~/home-lab/projects/<name>/<name>.container \
#                        ~/.config/containers/systemd/<name>.container
#                  systemctl --user daemon-reload && systemctl --user start <name>
[Unit]
Description=<name> — what it does
Wants=network-online.target
After=network-online.target

[Container]
ContainerName=<name>
Image=localhost:5000/<name>:0.1.0
Network=ai-net.network
EnvironmentFile=%h/home-lab/projects/<name>/<name>.env

Label=traefik.enable=true
# Label=traefik.http.services.<name>.loadbalancer.server.port=8080

[Service]
Restart=on-failure
TimeoutStartSec=300

[Install]
WantedBy=default.target
