# portainer/portainer.container
# Portainer CE — web UI for managing this rootless Podman instance.
# Quadlet unit (systemd --user). Symlinked into ~/.config/containers/systemd/.
[Unit]
Description=Portainer CE — container management UI
Wants=network-online.target podman.socket
After=network-online.target podman.socket

[Container]
ContainerName=portainer
# Pinned by digest (LTS, built 2026-06-03) — no clean version tag is published.
# To update: re-pull :lts, read its new digest, replace below, restart.
Image=docker.io/portainer/portainer-ce@sha256:0c5232e9a972ea2041070544e6afc4fb4cf1d21e1463ab37da3602de5360d086
Network=ai-net.network
# Podman API (rootless) — RW so Portainer can start/stop/inspect containers.
Volume=%t/podman/podman.sock:/var/run/docker.sock
# Persistent data.
Volume=portainer-data.volume:/data
# Auto-routed by Traefik at http://portainer.lab.lan (host derived from ContainerName).
# Port MUST be declared: the image exposes 8000/9000/9443, so Traefik can't guess.
Label=traefik.enable=true
Label=traefik.http.services.portainer.loadbalancer.server.port=9000

[Service]
Restart=on-failure
TimeoutStartSec=300

[Install]
WantedBy=default.target
